# systemd/certbot-restart.timer
[Unit]
Description=Timer to restart certbot service

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
